# project/exp-warning.cron
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
# Jalankan pengecekan peringatan expired setiap hari jam 9 pagi
0 9 * * * root /usr/bin/exp-warning
# Jalankan pengecekan peringatan expired setiap hari jam 18 sore (double check)
0 18 * * * root /usr/bin/exp-warning